# bsw/systemd/pki.service
[Unit]
Description=PKI Box Core Service
After=network.target
Wants=pki-hsm.service

[Service]
Type=simple
WorkingDirectory=/opt/pki-on-box
ExecStart=/opt/pki-on-box/venv/bin/python -m asw.PKI.main serve
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
User=pki
Group=pki

ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/pki /var/log/pki-box /opt/pki-on-box/data /opt/pki-on-box/logs
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
